# bridge.service
[Unit]
Description=Links two serial ports together

[Service]
ExecStart=/opt/serial-bridge/bridge.py

[Install]
WantedBy=local-fs.target
